# 98-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/mike/.Xauthority", RUN+="/home/mike/bin/monitor-hotplug.sh"
